# squashfs-root/etc/uci-defaults/luci-upnp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@UPnP[-1]
	add ucitrack UPnP
	set ucitrack.@UPnP[-1].init=upnp
EOF

rm -f /tmp/luci-indexcache
exit 0
