# squashfs-root/etc/uci-defaults/luci-dmz
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@DMZ[-1]
	add ucitrack DMZ
	set ucitrack.@DMZ[-1].init=DMZ
EOF

rm -f /tmp/luci-indexcache
exit 0
