# squashfs-root/etc/init.d/free_cache
#!/bin/sh /etc/rc.common

START=90
boot() {
	/bin/free_caches.sh &
}
