# squashfs-root/etc/init.d/lc
#!/bin/sh /etc/rc.common

######################################################################
#
# Function  : init script, lc
#
# Creater   : xiangda.zhang@feixun.com.cn, 2016-05-04
#
# Copyright : Shanghai Phicomm
#
######################################################################

START=80

boot() {
  /usr/bin/lc &
}

start() {
  /usr/bin/lc &
}

stop() {
  killall -9 lc
}

restart() {
  stop
  start
}
