# squashfs-root/etc/init.d/light_manage
#!/bin/sh /etc/rc.common

START=30

start() {
    light_manage &
}

stop() {
	killall -9 light_manage
}
