# squashfs-root/etc/init.d/luci_generatecache
#!/bin/sh /etc/rc.common

START=45
start() {
	/www/cgi-bin > /dev/null 2>1
}
