# squashfs-root/etc/init.d/update_up
#!/bin/sh /etc/rc.common

START=99

start() {
	[ -f "/etc/config/update_up" ] && {
		update_up permanent
		rm -rf "/etc/config/update_up"
	}
}
