# squashfs-root/etc/hotplug.d/iface/26-light
#!/bin/sh

if [ "$ACTION" = "ifdown" ]; then
    uci set network.wan.status="down"
elif [ "$ACTION" = "ifup" ]; then
	uci set network.wan.status="up"
fi


