# squashfs-root/etc/hotplug.d/iface/27-sysntpd
#!/bin/sh

[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
	/etc/init.d/sysntpd restart
}

