# squashfs-root/etc/hotplug.d/iface/55-parentctl
#!/bin/sh

if [ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ]; then
	/etc/init.d/parentctl restart
fi
