# squashfs-root/etc/hotplug.d/iface/71-appportfwd
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
	/etc/init.d/portfwd restart
fi
