# squashfs-root/etc/hotplug.d/iface/72-UPnP
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
        /etc/init.d/upnp restart
fi
