# squashfs-root/etc/hotplug.d/iface/73-dmz
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
	/etc/init.d/DMZ restart
fi
