# squashfs-root/etc/hotplug.d/rth/00-wan-link
#!/bin/sh
if [ "$EVENT" = "wan_link" ]; then
	if [ "$ACTION" = "up" ]; then
		logger "wan up!"
	elif [ "$ACTION" = "down" ]; then
		logger "wan down!"
	fi
fi
