# squashfs-root/etc/hotplug.d/rth/01-inet-link
#!/bin/sh
if [ "$EVENT" = "inet_link" ]; then
	if [ "$ACTION" = "up" ]; then
		logger "inet up!"
	elif [ "$ACTION" = "down" ]; then
		logger "inet down!"
	fi
fi
