# squashfs-root/usr/lib/opkg/info/button.list
/etc/nvram2uci/button
/etc/init.d/button
/usr/bin/button
/etc/config/button
