# squashfs-root/usr/lib/opkg/info/capp.list
/usr/bin/capp
/etc/init.d/capp
/etc/uci-defaults/99-capp
/usr/share/capp/firewall.include
