# squashfs-root/usr/lib/opkg/info/config_util.list
/usr/bin/config_util
/etc/init.d/config_util
