# squashfs-root/usr/lib/opkg/info/encryconfig.list
/usr/bin/encryconfig
/etc/init.d/genrandkey
/usr/bin/genrandkey
