# squashfs-root/usr/lib/opkg/info/ip_monitor.list
/etc/init.d/ip_monitor
/usr/bin/ip_monitor
