# squashfs-root/usr/lib/opkg/info/luci-mod-base.list
/usr/lib/lua/luci/view/h5/wifiConfig.htm
/usr/lib/lua/luci/controller/admin/page.lua
/www/luci-static/images/topology.png
/etc/config/luci
/usr/lib/lua/luci/phicomm/validator.lua
/www/luci-static/images/button_disable.png
/www/luci-static/images/signal_4.png
/www/luci-static/images/wisp_h5.png
/etc/init.d/luci_fixtime
/usr/lib/lua/luci/view/pc/netState.htm
/usr/lib/lua/luci/controller/admin/protest.lua
/www/luci-static/images/phicomm.png
/www/luci-static/images/signal.gif
/www/luci-static/dynaform/up.css
/www/luci-static/images/pic_code.png
/www/luci-static/images/signal_lock_1.png
/usr/lib/lua/luci/view/h5/setLgPwd.htm
/usr/lib/lua/luci/phicomm/setting.lua
/www/luci-static/images/wifi.png
/www/luci-static/dynaform/icon_fonts_style.css
/www/luci-static/images/wisp_pc.png
/www/luci-static/dynaform/frame.css
/www/luci-static/dynaform/h5-en-us.css
/www/luci-static/images/color_line.png
/www/luci-static/images/upgrade_rocket.png
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/view/h5/guide.htm
/usr/lib/lua/luci/view/error404.htm
/usr/lib/lua/luci/controller/admin/dev_feature.lua
/www/luci-static/images/signal_lock_2.png
/usr/lib/lua/luci/controller/ds.lua
/usr/lib/lua/luci/view/h5/menu.htm
/www/luci-static/images/pic-null.png
/usr/lib/lua/luci/data/guide.lua
/www/luci-static/dynaform/upgrade.js
/www/luci-static/images/dog.png
/www/luci-static/lib/json.min.js
/usr/lib/lua/luci/view/pc/guideWifiSet.htm
/usr/lib/lua/luci/phicomm/debug.lua
/usr/lib/lua/luci/view/h5/guideWifiSet.htm
/www/luci-static/dynaform/routing.js
/www/luci-static/images/network.png
/www/luci-static/dynaform/pc-en-us.css
/www/luci-static/lib/libs.js
/www/luci-static/images/router_blue.png
/usr/lib/lua/luci/view/pc/setLgPwd.htm
/www/luci-static/lib/flexible.min.js
/www/luci-static/images/signal_lock_3.png
/usr/lib/lua/luci/view/pc/userAgreement.htm
/etc/cfgAdaptor/system.sh
/usr/lib/lua/luci/controller/admin/wireless.lua
/www/luci-static/dynaform/vendor.js
/usr/lib/lua/luci/view/pc/guide.htm
/www/luci-static/images/pic_router.png
/www/luci-static/images/pic_left.png
/www/luci-static/images/pic-router_h5.png
/www/luci-static/dynaform/verify.js
/www/luci-static/images/pic-router_pc.png
/www/luci-static/dynaform/h5.js
/www/luci-static/dynaform/zh-cn.js
/usr/lib/lua/luci/phicomm/util.lua
/usr/lib/lua/luci/controller/admin/device.lua
/www/luci-static/dynaform/error.js
/usr/lib/lua/luci/data/page.lua
/www/luci-static/images/signal_lock_4.png
/www/luci-static/images/button.png
/usr/lib/lua/luci/phicomm/statistic.lua
/www/luci-static/lib/stateman.min.js
/www/luci-static/images/circle.png
/www/luci-static/lib/respond.min.js
/usr/lib/lua/luci/controller/admin/find_pppoe_pwd.lua
/usr/lib/lua/luci/controller/admin/device_mng.lua
/www/luci-static/dynaform/frame.js
/www/luci-static/images/app-icon/reboot.png
/usr/lib/lua/luci/phicomm/fs.lua
/etc/config/timereboot
/etc/init.d/luci_generatecache
/etc/cfgAdaptor/luci.sh
/www/luci-static/images/loading.gif
/usr/lib/lua/luci/data/wireless.lua
/usr/lib/lua/luci/view/userAgreementContent.htm
/www/luci-static/images/password.png
/usr/lib/lua/luci/view/index.htm
/usr/lib/lua/luci/view/pc/wifiConfig.htm
/www/luci-static/fonts/icofont.eot
/usr/lib/lua/luci/data/device_mng.lua
/www/luci-static/dynaform/h5.css
/www/luci-static/images/router-white.png
/usr/lib/lua/luci/dispatcher.lua
/usr/lib/lua/luci/view/h5/checkRouterLink.htm
/www/luci-static/images/router.png
/usr/lib/lua/luci/view/h5/editClient.htm
/www/luci-static/dynaform/pc.css
/usr/lib/lua/luci/view/h5/userAgreement.htm
/usr/lib/lua/luci/view/h5/login.htm
/www/luci-static/dynaform/up.js
/www/luci-static/lib/jquery.min.js
/www/luci-static/images/signal_1.png
/usr/lib/lua/luci/view/pc/checkRouterLink.htm
/www/luci-static/dynaform/en-us.js
/www/luci-static/images/pic_computer.png
/www/luci-static/fonts/icofont.ttf
/usr/lib/lua/luci/view/pc/menu.htm
/usr/lib/lua/luci/view/pc/deviceManage.htm
/usr/lib/lua/luci/view/h5/deviceManage.htm
/www/luci-static/images/phiwifi.png
/www/luci-static/dynaform/style_1080p.css
/www/luci-static/lib/ajax.js
/www/luci-static/images/bg_map.jpg
/usr/lib/lua/luci/phicomm/lib/shell.lua
/usr/lib/lua/luci/view/pc/findPppoePwd.htm
/usr/lib/lua/luci/view/h5/findPppoePwd.htm
/usr/lib/lua/luci/view/pc/App.htm
/www/luci-static/images/signal_2.png
/www/luci-static/dynaform/pc.js
/www/luci-static/dynaform/ie.css
/usr/lib/lua/luci/view/pc/upgradeWgt.htm
/usr/lib/lua/luci/view/h5/netState.htm
/usr/lib/lua/luci/phicomm/error.lua
/usr/lib/lua/luci/view/pc/login.htm
/usr/lib/lua/luci/view/pc/modifyPwd.htm
/usr/lib/lua/luci/view/error500.htm
/usr/lib/lua/luci/phicomm/lib/bit.lua
/usr/lib/lua/luci/view/pc/timeReboot.htm
/www/luci-static/images/bg_botton.png
/www/luci-static/images/app-icon/lock.png
/www/luci-static/images/signal_3.png
/www/luci-static/images/favicon.ico
/usr/lib/lua/luci/view/h5/modifyPwd.htm
/www/luci-static/images/pic_bottombar.png
/usr/lib/lua/luci/view/h5/App.htm
