# squashfs-root/usr/lib/opkg/info/luci-app-parent-ctrl.list
/usr/lib/lua/luci/view/h5/parentCtrlRule.htm
/etc/hotplug.d/iface/55-parentctl
/usr/lib/lua/luci/view/h5/parentCtrlAddRule.htm
/usr/lib/lua/luci/view/pc/parentCtrl.htm
/usr/lib/lua/luci/view/h5/devList.htm
/etc/config/parentctl
/usr/lib/lua/luci/controller/admin/parent_ctrl.lua
/www/luci-static/images/app-icon/parent_ctrl.png
/usr/lib/lua/luci/data/parent_ctrl.lua
/usr/lib/lua/luci/view/pc/timePick.htm
/usr/lib/lua/luci/view/h5/parentCtrl.htm
/etc/init.d/parentctl
