# squashfs-root/usr/lib/opkg/info/oui_update.list
/usr/bin/oui_update
/etc/init.d/oui_sync
