# squashfs-root/usr/lib/opkg/info/rp-pppoe-server.list
/etc/ppp/pppoe-server-options
/etc/init.d/pppoe-server
/etc/ppp/pap-secrets
/etc/ppp/plugins/rp-pppoe.so
/etc/default/pppoe-server
/usr/sbin/pppoe-server
