# squashfs-root/usr/lib/opkg/info/sys_protect_daemon.list
/etc/init.d/sys_protect_daemon
/usr/bin/sys_protect_daemon
