# squashfs-root/usr/lib/opkg/info/telnetd_startup.list
/usr/bin/telnetd_startup
/etc/init.d/telnetd_startup
