# squashfs-root/usr/lib/opkg/info/tnat.list
/usr/bin/lc
/etc/hotplug.d/iface/80-tnat
/etc/init.d/lc
/etc/config/tnat
