# squashfs-root/usr/lib/opkg/info/update_usrproto.list
/etc/init.d/update_up
/usr/sbin/update_up
/etc/config/update_up
