# squashfs-root/etc/hotplug.d/ip_event/99-adportal
#!/bin/sh

[ "$ACTION" != "OFFLINE" ] && return

ipset del ad_white "$MAC" >/dev/null 2>&1
